# crontab
5 4 * * * python3 -W ignore ~/meteostat-server/tasks/clean.py
5 5 * * * sh ~/meteostat-server/tasks/db_import.sh
